# com.newsresearch.memory_os.compact.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Label must match the filename (without .plist) -->
    <key>Label</key>
    <string>com.newsresearch.memory_os.compact</string>

    <!-- Run once per day at 03:00 local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Command: project venv python + python -m memory_os compact -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/oleksii/Documents/News Research Automation/venv_auto/bin/python</string>
        <string>-m</string>
        <string>memory_os</string>
        <string>compact</string>
    </array>

    <!-- Working directory must be the project root so imports resolve -->
    <key>WorkingDirectory</key>
    <string>/Users/oleksii/Documents/News Research Automation</string>

    <!-- Logs -->
    <key>StandardOutPath</key>
    <string>/Users/oleksii/Documents/News Research Automation/logs/memory_os_compact.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/oleksii/Documents/News Research Automation/logs/memory_os_compact.err</string>

    <!-- Load at login, not at boot (user agent) -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Retry on failure after 60s -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
